# backend/deploy/cantina.service
[Unit]
Description=Cantina household grocery + meal tracker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/penguin/projects/cantina/backend/src
ExecStart=/home/penguin/projects/cantina/backend/src/venv/bin/uvicorn api:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=3
# Keep the process tightly scoped -- it only ever needs its own files.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/penguin/projects/cantina/backend/src

[Install]
WantedBy=default.target
